# legacy/nextcloud-server/config/www-data.cron
* * * * * php -f /var/www/html/cron.php

# @reboot   php /var/www/html/occ files:scan --all
# * * * * * php /var/www/html/occ files:scan --all
