# legacy/nextcloud-cron/config/tasks.cron
@reboot   chmod -R 777 /var/www/html/data && php occ files:scan /var/www/html/data
* * * * * chmod -R 777 /var/www/html/data && php occ files:scan /var/www/html/data
