# services/nextcloud-server/config/root.cron
@reboot   sudo -u www-data php /var/www/html/occ files:scan --all && sudo chmod -R 777 /var/www/html/data
* * * * * sudo -u www-data php /var/www/html/occ files:scan --all && sudo chmod -R 777 /var/www/html/data
